# profiledir/airootfs/etc/systemd/system/anonos-firstboot.service
[Unit]
Description=AnonOS First Boot System Setup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/anonos-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
